# examples/buildroot-polyglot-squashfs/assets/systemd/rust-hello.service
[Unit]
Description=Rust hello demo service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/rust-hello
Restart=on-failure

[Install]
WantedBy=multi-user.target
