# examples/imported-rootfs-polyglot-git/assets/systemd/polyglot.service
[Unit]
Description=Gaia Polyglot Git Smoke
After=network.target

[Service]
EnvironmentFile=/etc/default/runtime.env
ExecStart=/usr/bin/rust-app
Restart=on-failure

[Install]
WantedBy=multi-user.target
